# root/etc/s6-overlay/s6-rc.d/init-ubooquity-config/up
/etc/s6-overlay/s6-rc.d/init-ubooquity-config/run
